# server/systemd/mac-mount-check.service
[Unit]
Description=fstack — check and remount stale Mac SSHFS mounts
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-mac-mounts
