# server/systemd/moto-cleanup.service
[Unit]
Description=fstack — reap orphan dev-server and chrome processes
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup-stale
