# server/systemd/node-modules-gc.timer
[Unit]
Description=fstack — nightly node_modules + orphan cleanup

[Timer]
OnCalendar=*-*-* 04:12:00
Persistent=true
Unit=node-modules-gc.service

[Install]
WantedBy=timers.target
