# server/systemd/node-modules-gc.service
[Unit]
Description=fstack — garbage-collect stale node_modules directories
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/node-modules-gc
ExecStartPost=/usr/local/bin/kill-claude-orphans
